# roles/kubernetes/files/etc/systemd/system/datadog-agent.service.d/10-kubernetes.conf
[Service]
# force dd-agent's tagger to talk to kubelet for tags
# Workaround for https://github.com/DataDog/datadog-agent/pull/7583
Environment=KUBERNETES=true